# systemd/apache-age/apache-age.container
[Unit]
Description=Apache AGE — openCypher graph queries as a PostgreSQL extension
Documentation=https://age.apache.org/getstarted/quickstart
Wants=network-online.target apache-age-data-volume.service
After=network-online.target apache-age-data-volume.service

[Container]
AutoUpdate=registry
Image=docker.io/apache/age:latest
ContainerName=systemd-apache-age
PublishPort=5455:5432
Volume=apache-age-data:/var/lib/postgresql/data
Environment=POSTGRES_USER=age_user
Environment=POSTGRES_DB=age_db
Secret=apache-age-password,type=env,target=POSTGRES_PASSWORD
HealthCmd=["pg_isready","-U","age_user","-d","age_db"]
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=30s
HealthTimeout=10s
Label=app=apache-age

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target default.target
